# pi-config/systemd/alohamini-cam.service
[Unit]
Description=AlohaMini Camera Server
After=alohamini-host.service
Wants=alohamini-host.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/lerobot_alohamini
Environment=CAMS=2,4,6,8
ExecStart=/home/pi/miniforge3/envs/lerobot/bin/python examples/alohamini/robot_cam_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alohamini-cam

[Install]
WantedBy=multi-user.target
